# cron.d/ctdb_monitor.sh
# Every 3 minutes, ctdb monitor are running
*/3 * * * * root /etc/ctdb/monitor_ctdb.sh
