# cron.d/ctdb_monitor
# Every minute, monitor ctdb are running
* * * * * root /etc/ctdb/cron-seconds
